# roles/fix_ubunu_zfs_encrypted_swap_setup/templates/hibernate-configure.service.j2
[Unit]
Description=Hibernate Configuration Service
After=swap.target

[Service]
Type=oneshot
ExecStart=/root/hibernate_configure.sh
RemainAfterExit=false

[Install]
WantedBy=multi-user.target
