# roles/power_management/templates/force-hibernate.timer.j2
[Unit]
Description=Run Force Hibernate script every 20 seconds

[Timer]
OnUnitActiveSec=20s
AccuracySec=5s
Persistent=true

[Install]
WantedBy=timers.target
