# roles/power_management/templates/force-hibernate.service.j2
[Unit]
Description=Force Hibernate on Critical Battery

[Service]
Type=oneshot
ExecStart=/usr/local/bin/force-hibernate.sh
